# systemd/scope1090-persist.timer
[Unit]
Description=scope1090 — hourly DB persistence
Requires=scope1090-persist.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
